# crontab
0 * * * * cd /app && poetry run poe welcome > /var/log/cron.log 2>&1
