# crontab
0 * * * * cd /app && poetry run poe index > /var/log/cron.log 2>&1
